# .github/deploy/hobot.service
[Unit]
Description=Hobot FastAPI Service
After=network.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user/hobot-service/hobot
Environment="PATH=/home/ec2-user/hobot-service/hobot/venv/bin"
# Python 모듈로 실행하여 가상환경의 Python을 사용
ExecStart=/home/ec2-user/hobot-service/hobot/venv/bin/python -m gunicorn -c gunicorn.conf.py asgi:asgi_app
Restart=always
RestartSec=10

# 로그 디렉토리 생성 및 권한 설정
ExecStartPre=/bin/mkdir -p /home/ec2-user/hobot-service/hobot/logs
ExecStartPre=/bin/chmod 755 /home/ec2-user/hobot-service/hobot/logs

# 로그 설정 (journal로 먼저 설정, 파일 리다이렉션은 gunicorn.conf.py에서 처리)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hobot

[Install]
WantedBy=multi-user.target